# docs/launchd/com.codex.relay.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.codex.relay</string>

    <key>ProgramArguments</key>
    <array>
      <string>/Users/cypher-server/Documents/code/codex-flow/.venv/bin/python</string>
      <string>-m</string>
      <string>codex_relay</string>
      <string>daemon</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/cypher-server/Documents/code/codex-flow/logs/codex-relay.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/cypher-server/Documents/code/codex-flow/logs/codex-relay.err.log</string>

    <key>WorkingDirectory</key>
    <string>/Users/cypher-server/Documents/code/codex-flow</string>

    <key>EnvironmentVariables</key>
    <dict>
      <key>PATH</key>
      <string>/Users/cypher-server/Documents/code/codex-flow/.venv/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
  </dict>
</plist>
